# systemd/xray-mem-guard.service
[Unit]
Description=Restart the Xray core when memory usage is high (3X-UI mem guard)
Documentation=https://github.com/<your-github>/xray-mem-guard
After=x-ui.service

[Service]
Type=oneshot
Nice=10
ExecStart=/usr/local/bin/xray-mem-guard.sh
